# src/websocket.conf
[program:websocket]
command=php /var/www/tetrinet/src/server/artisan websockets:serve --port 10000
directory=/var/www/tetrinet/src/server
user=www-data
autostart=true
autorestart=true